# local.Yosemite.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.Yosemite</string>
    <key>WorkingDirectory</key>
    <string>/Users/Parent/Python/yosemite-camping-master</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/Parent/Python/yosemite-camping-master/campsites.sh</string>
        <string>2018-03-01</string>
        <string>2018-03-03</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardErrorPath</key>
    <string>/tmp/local.Yosemite.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/local.Yosemite.out</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
